# infrastructure/scripts/nse-api-staging.service
# systemd service — FastAPI (STAGING)
#
# Uses stg_ DynamoDB prefix, reads secrets from /nse/staging/ in SSM.
# Runs on port 9001 alongside prod (port 9000) on the same EC2.
#
# Install:
#   sudo cp infrastructure/scripts/nse-api-staging.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable nse-api-staging
#   sudo systemctl start nse-api-staging
#
# Logs:   sudo journalctl -u nse-api-staging -f
# Reload: sudo systemctl restart nse-api-staging

[Unit]
Description=NSE Stock Dashboard — FastAPI (staging)
After=network.target
Wants=network.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nse-staging/backend
Environment=PATH=/opt/nse-staging/venv/bin:/usr/bin:/bin
Environment=STAGE=staging
EnvironmentFile=-/opt/nse-staging/.env

ExecStart=/opt/nse-staging/venv/bin/uvicorn \
    app.main:app \
    --host 127.0.0.1 \
    --port 9001 \
    --workers 1 \
    --log-level info \
    --no-access-log

Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

MemoryLimit=600M
CPUQuota=80%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nse-api-staging

[Install]
WantedBy=multi-user.target
