# infrastructure/scripts/nse-worker-dev.service
# systemd service — Playwright scraping worker (DEV stage)
#
# Install:
#   sudo cp infrastructure/scripts/nse-worker-dev.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable nse-worker-dev
#   sudo systemctl start nse-worker-dev
#
# Logs:   sudo journalctl -u nse-worker-dev -f

[Unit]
Description=NSE Stock Dashboard — Playwright Scraping Worker (dev)
After=network.target nse-api-dev.service
Wants=network.target

[Service]
Type=exec
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nse-dev/backend
Environment=PATH=/opt/nse-dev/venv/bin:/usr/bin:/bin
Environment=PLAYWRIGHT_BROWSERS_PATH=/home/ubuntu/.cache/ms-playwright
Environment=STAGE=dev
EnvironmentFile=-/opt/nse-dev/backend/.env

ExecStart=/opt/nse-dev/venv/bin/python -m app.worker

Restart=always
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5

MemoryLimit=500M
CPUQuota=70%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nse-worker-dev

[Install]
WantedBy=multi-user.target
